# systemd/mobileshell.service
[Unit]
Description=MobileShell {{USER}} - Web-based shell access
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory=/home/{{USER}}
ExecStart=/opt/{{USER}}-mobileshell run
StateDirectory=mobileshell-{{USER}}
StateDirectoryMode=0700
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
